# crontab
* * * * * /app/scripts/backup.sh >> /app/logs/backup.log 2>&1
